# wallpaper-randomizer/wallpaper-randomizer.service
[Unit]
Description=Wallpaper Randomizer

[Service]
ExecStart=/home/stasnocap/.local/scripts/wallpaper-randomizer

[Install]
WantedBy=default.target
